# deployment/datapusher-uwsgi.conf
[program:datapusher-uwsgi]

command=/usr/lib/ckan/datapusher/bin/uwsgi -i /usr/lib/ckan/datapusher/src/datapusher/deployment/datapusher-uwsgi.ini
environment=PATH="/usr/lib/ckan/datapusher/bin:%(ENV_PATH)s"

; Start just a single worker for the uwsgi master process.
; If you need more concurrency, adjust workers and threads
; in datapusher-uwsgi.ini
numprocs=1
process_name=%(program_name)s-%(process_num)02d

; Log files - change this to point to the same location as
; the existing CKAN log files
stdout_logfile=/etc/ckan/default/datapusher.OUT
stderr_logfile=/etc/ckan/default/datapusher.ERR

; Make sure that the worker is started on system start and automatically
; restarted if it crashes unexpectedly.
autostart=true
autorestart=true

; Number of seconds the process has to run before it is considered to have
; started successfully.
startsecs=10

; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs = 600

; Required for uWSGI as it does not obey SIGTERM.
stopsignal=QUIT
